# infrastructure/supervisor/supervisord.conf

[inet_http_server]
port=*:9001  
username=user 
password=pass  

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
user=root

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:symfony-messenger-consumer]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
command= php bin/console messenger:consume async --time-limit=3600 --memory-limit=128M --no-interaction
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/symfony-messenger-consumer.log


